# scripts/dicton.service
[Unit]
Description=Dicton Voice Transcription
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
Environment=DISPLAY=:0
WorkingDirectory=__PROJECT_DIR__
ExecStart=__PROJECT_DIR__/venv/bin/dicton
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
